# kali-dpkg-packages/auto-vpn/DEBIAN/postinst
#!/bin/bash
set -e

# Enable and start the openvpn service
systemctl enable auto-vpn.service
systemctl start auto-vpn.service