# kali-dpkg-packages/mount-kali-dpkg-packages/DEBIAN/postinst
#!/bin/bash
set -e

# Enable and start the mount-kali-dpkg-packages service
systemctl enable mount-kali-dpkg-packages.service
systemctl start mount-kali-dpkg-packages.service
